# deploy/bookcabinet-ui.service
[Unit]
Description=BookCabinet Node.js UI (Express + React)
After=network.target pigpiod.service bookcabinet-daemon.service
Wants=pigpiod.service bookcabinet-daemon.service

[Service]
Type=simple
User=admin42
Group=admin42
WorkingDirectory=/home/admin42/bookcabinet
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Environment
Environment=NODE_ENV=production
Environment=PORT=5000
# Киоск локальный (chromium на этом же Pi). Раньше HOST игнорировался кодом
# и сервер слушал 127.0.0.1; после фикса index.ts значение действует —
# 0.0.0.0 открыл бы API всей сети без авторизации (синглтон-сессия).
Environment=HOST=127.0.0.1
Environment=BOOKCABINET_ROOT=/home/admin42/bookcabinet
Environment=PYTHON_BRIDGE_PATH=/home/admin42/bookcabinet/bookcabinet/bridge.py
Environment=AUTH_DAEMON_PATH=/home/admin42/bookcabinet/bookcabinet/services/auth_shutter_daemon.py
Environment=MOCK_MODE=false

# Watchdog
WatchdogSec=60
NotifyAccess=all

[Install]
WantedBy=multi-user.target
