# deploy/bookcabinet-api.service
# BookCabinet — целевой python-бэкенд (aiohttp, решение 2026-06-12).
#
# На шкафу боевой юнит называется bookcabinet.service и УЖЕ запускает python.
# Этот файл — эталонное содержимое. Деплой = обновить /etc/systemd/system/bookcabinet.service
# этим содержимым (или установить как bookcabinet-api и отключить старый), см. deploy/DEPLOY.md.
# Перед этим: npm run build (aiohttp отдаёт dist/public) и deploy/migrate.sh.
#
# Сверено с устройством 2026-06-13: pcscd через ПОСТОЯННЫЙ pcscd-daemon.service
# (Type=forking, sleep 5 для USB, держит claim ACR заранее → нет гонки старта NFC).
# НЕ через pcscd.service/.socket: их --auto-exit + параллельный daemon давали дубль
# pcscd → LIBUSB_BUSY и NFC не поднимался. socket-механизм на устройстве отключён.
# HOST локальный (раньше 0.0.0.0 — открывало API в сеть при синглтон-авторизации).
[Unit]
Description=BookCabinet API (Python aiohttp)
After=network.target pigpiod.service pcscd-daemon.service
Wants=pigpiod.service pcscd-daemon.service

[Service]
Type=simple
User=admin42
WorkingDirectory=/home/admin42/bookcabinet
ExecStart=/usr/bin/python3 -m bookcabinet.main
Restart=always
RestartSec=3
# Киоск локальный: наружу не слушаем (авторизация — синглтон, см. docs/STATE.md)
Environment=HOST=127.0.0.1
Environment=PORT=5000
Environment=MOCK_MODE=false
Environment=IRBIS_MOCK=false
Environment=DATABASE_PATH=/home/admin42/bookcabinet/bookcabinet/data/shelf_data.db
Environment=PYTHONUNBUFFERED=1
SupplementaryGroups=dialout gpio i2c spi plugdev

[Install]
WantedBy=multi-user.target
